# deploy/autoresearch-trigger.service
[Unit]
Description=Elastifund Autoresearch Trigger
After=network.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/polymarket-trading-bot
EnvironmentFile=-/home/ubuntu/polymarket-trading-bot/config/btc5_strategy.env
EnvironmentFile=-/home/ubuntu/polymarket-trading-bot/state/btc5_capital_stage.env
EnvironmentFile=/home/ubuntu/polymarket-trading-bot/.env
ExecStart=/usr/bin/python3 scripts/autoresearch_trigger.py
StandardOutput=journal
StandardError=journal
