# deploy/btc-5min-maker.service
[Unit]
Description=BTC 5-Minute Maker Bot (Instance 2)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/polymarket-trading-bot
Environment="PYTHONPATH=/home/ubuntu/polymarket-trading-bot:/home/ubuntu/polymarket-trading-bot/bot:/home/ubuntu/polymarket-trading-bot/polymarket-bot" "BTC5_MAX_ABS_DELTA=0.00015" "BTC5_UP_MAX_BUY_PRICE=0.51" "BTC5_DOWN_MAX_BUY_PRICE=0.51" "BTC5_ENABLE_RECENT_REGIME_SKEW=1" "BTC5_RECENT_REGIME_FILLS=12" "BTC5_REGIME_MIN_FILLS_PER_DIRECTION=5" "BTC5_REGIME_MIN_PNL_GAP_USD=20" "BTC5_REGIME_WEAKER_DIRECTION_QUOTE_TICKS=0"
EnvironmentFile=/home/ubuntu/polymarket-trading-bot/.env
ExecStart=/usr/bin/python3 bot/btc_5min_maker.py --continuous --live --verbose
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=btc5maker

[Install]
WantedBy=multi-user.target
