# deploy/btc5-autoresearch.timer
[Unit]
Description=Run BTC5 Autoresearch every 3 hours

[Timer]
OnBootSec=15min
OnUnitActiveSec=3h
Persistent=true
Unit=btc5-autoresearch.service

[Install]
WantedBy=timers.target
