# deploy/btc5-dual-autoresearch-morning.service
[Unit]
Description=BTC5 Dual Autoresearch Morning Report
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/polymarket-trading-bot
ExecStart=/bin/bash -lc 'mkdir -p state reports/autoresearch && export PYTHONPATH=/home/ubuntu/polymarket-trading-bot:/home/ubuntu/polymarket-trading-bot/bot:/home/ubuntu/polymarket-trading-bot/polymarket-bot && /usr/bin/python3 scripts/btc5_dual_autoresearch_ops.py morning-report --window-hours 24'
TimeoutStartSec=300
StandardOutput=journal
StandardError=journal
SyslogIdentifier=btc5dualautoresearchmorning
