# deploy/btc5-market-model-autoresearch.service
[Unit]
Description=BTC5 Market-Model Mutation Cycle
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/polymarket-trading-bot
ExecStart=/bin/bash -lc 'mkdir -p state reports/autoresearch && export PYTHONPATH=/home/ubuntu/polymarket-trading-bot:/home/ubuntu/polymarket-trading-bot/bot:/home/ubuntu/polymarket-trading-bot/polymarket-bot && /usr/bin/python3 scripts/run_btc5_market_model_mutation_cycle.py'
Restart=on-failure
RestartSec=60
TimeoutStartSec=1800
StandardOutput=journal
StandardError=journal
SyslogIdentifier=btc5marketautoresearch
