# deploy/execution-liveness.service
[Unit]
Description=Execution liveness monitor — detect trading stalls
After=network-online.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/polymarket-trading-bot
Environment="PYTHONPATH=/home/ubuntu/polymarket-trading-bot:/home/ubuntu/polymarket-trading-bot/bot:/home/ubuntu/polymarket-trading-bot/polymarket-bot"
ExecStart=/usr/bin/python3 -m bot.execution_liveness --once
StandardOutput=journal
StandardError=journal
SyslogIdentifier=execution-liveness
