# deploy/xrp-5min-maker.service
[Unit]
Description=XRP 5-Minute Maker Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/polymarket-trading-bot
Environment="PYTHONPATH=/home/ubuntu/polymarket-trading-bot:/home/ubuntu/polymarket-trading-bot/bot:/home/ubuntu/polymarket-trading-bot/polymarket-bot" "BTC5_ENABLE_RECENT_REGIME_SKEW=1" "BTC5_RECENT_REGIME_FILLS=12" "BTC5_REGIME_MIN_FILLS_PER_DIRECTION=5" "BTC5_REGIME_MIN_PNL_GAP_USD=20" "BTC5_REGIME_WEAKER_DIRECTION_QUOTE_TICKS=0" "BTC5_ENABLE_PROBE_AFTER_DAILY_LOSS=1" "BTC5_ENABLE_PROBE_AFTER_RECENT_LOSS=1" "BTC5_PROBE_RECENT_FILLS=8" "BTC5_PROBE_RECENT_MIN_PNL_USD=0.0" "BTC5_PROBE_QUOTE_TICKS=0" "BTC5_RETRY_POST_ONLY_CROSS=1" "BTC5_RETRY_POST_ONLY_SAFETY_TICKS=1" "BTC5_ASSET_SLUG_PREFIX=xrp" "BTC5_ASSET_BINANCE_SYMBOL=XRPUSDT" "BTC5_BINANCE_WS_URL=wss://stream.binance.com:9443/ws/xrpusdt@trade" "BTC5_BINANCE_TICKER_URL=https://api.binance.com/api/v3/ticker/price?symbol=XRPUSDT" "BTC5_DB_PATH=data/xrp_5min_maker.db"
EnvironmentFile=/home/ubuntu/polymarket-trading-bot/.env
EnvironmentFile=-/home/ubuntu/polymarket-trading-bot/config/btc5_strategy.env
EnvironmentFile=-/home/ubuntu/polymarket-trading-bot/state/btc5_capital_stage.env
EnvironmentFile=-/home/ubuntu/polymarket-trading-bot/config/xrp5_strategy.env
ExecStart=/usr/bin/python3 bot/btc_5min_maker.py --continuous --live
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=xrp5maker

[Install]
WantedBy=multi-user.target
